# deploy/research@.service
# systemd template unit. Instance specifier is the username to run as.
#
# Place at /etc/systemd/system/research@.service, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now research@<username>.service
#   sudo journalctl -u research@<username>.service -f
#
# Adjust the --host argument and --acme-email before enabling. ACME HTTP-01
# binds :80 and TLS binds :443, so the unit grants CAP_NET_BIND_SERVICE.

[Unit]
Description=Programming languages research corpus (running as %i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/research
ExecStart=/usr/bin/node source/research.js \
    --host research.example.com \
    --port 443 \
    --acme-email admin@example.com \
    --allowed-host research.example.com
Restart=on-failure
RestartSec=5
User=%i
Environment=NODE_ENV=production

AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/srv/research
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
